# station/autostart/start_reset_radio.desktop
[Desktop Entry]
Type=Application
Name=Streaming Station Controller
Exec=/home/pi/start_reset_radio
#Exec=/usr/bin/lxterminal
#Terminal=false
