# station/autostart/start_radio.desktop
[Desktop Entry]
Type=Application
Name=Radio Streaming
Exec=/home/pi/start_radio
#Exec=/usr/bin/lxterminal
#Terminal=false
